# debian/heartbeat-ldirectord.files
/etc/ha.d/resource.d/ldirectord
/etc/init.d/heartbeat-ldirectord
usr/share/doc/heartbeat-ldirectord/ldirectord.cf
etc/logrotate.d/ldirectord
usr/sbin/ldirectord
usr/sbin/supervise-ldirectord-config
usr/share/man/man8/ldirectord.8.gz
usr/share/man/man8/supervise-ldirectord-config.8.gz
