# debian/heartbeat-2.preinst
#!/bin/sh
#
# License: GNU General Public License (GPL)
#
# postinst script for heartbeat
#
# see: dh_installdeb(1)

if [ -x /etc/init.d/heartbeat ]; then
	if which invoke-rc.d >/dev/null 2>&1; then
		invoke-rc.d heartbeat stop
	else
		/etc/init.d/heartbeat stop
	fi
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
