# debian/ldirectord-2.files
/etc/ha.d/resource.d/ldirectord
/etc/init.d/ldirectord
/etc/logrotate.d/ldirectord
usr/sbin/ldirectord
usr/sbin/supervise-ldirectord-config
usr/share/doc/ldirectord-2/ldirectord.cf
usr/share/man/man8/ldirectord.8.gz
usr/share/man/man8/supervise-ldirectord-config.8.gz
