# debian/ldirectord.conffiles
/etc/logrotate.d/ldirectord
/etc/init.d/ldirectord
